# docker/php/supervisord.conf
[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0777

[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

; ==========================
; PHP-FPM
; ==========================
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
user=root
stdout_logfile=/var/www/storage/logs/php-fpm.log
redirect_stderr=true
stopasgroup=true
killasgroup=true

; ==========================
; Laravel Queue Worker
; ==========================
; Ordem de prioridade: cliente (KYT) → high (alerts) → default → import → ...
; --timeout=600 mata jobs após 10min (evita workers pendurados)
; numprocs=3 permite processar 3 jobs em paralelo
;
; Worker manual (se supervisor falhar):
;   docker exec keepcomply-QA-php php /var/www/artisan queue:work \
;     --queue=cliente,high,default,import,policyChanges,policy \
;     --tries=20 --sleep=3 --timeout=600
;
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=cliente,high,default,import,policyChanges,policy --tries=20 --sleep=3 --timeout=600
autostart=true
autorestart=true
user=www-data
numprocs=3
startsecs=3
startretries=10
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopasgroup=true
killasgroup=true
